# deploy/host/files/etc/systemd/system/pw-tmux.service.d/runtime.conf
# GOA: restore the .NET deploy toolchain (dotnet symlink + libicu + smbclient) that the
# canonical image does not bake, on every (re)start. Appends to the ExecStartPost list.
[Service]
ExecStartPost=-/usr/local/sbin/pw-ensure-dotnet.sh pw-tmux
