# systemd/claude-code-update.service
[Unit]
Description=Update Claude Code CLI
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/update-claude-code
