# systemd/project-terminal@.service
[Unit]
Description=ttyd terminal for project %i
After=network-online.target
Wants=network-online.target
[Service]
Type=simple
User=admin
Group=admin
Environment=HOME=/home/admin
ExecStart=/usr/local/bin/project-terminal-start %i
Restart=on-failure
RestartSec=3
[Install]
WantedBy=multi-user.target
