# systemd/project-terminal@.service
[Unit]
Description=ttyd terminal for project %i
After=network-online.target
Wants=network-online.target
# Let tmux sessions be restored before this terminal attaches, so
# `tmux new-session -A` reattaches to the restored session instead of creating
# an empty one. Wants= (not Requires=) so a restore failure never blocks the
# terminal — it just comes up fresh, as it did before persistence existed.
After=pw-tmux-persist.service
Wants=pw-tmux-persist.service
[Service]
Type=simple
User=admin
Group=admin
Environment=HOME=/home/admin
ExecStart=/usr/local/bin/project-terminal-start %i
Restart=on-failure
RestartSec=3
[Install]
WantedBy=multi-user.target
