# systemd/project-setup-terminal.service
[Unit]
Description=ttyd setup/auth terminal for Project Workbench Setup Wizard
After=network-online.target
Wants=network-online.target
[Service]
Type=simple
User=admin
Group=admin
Environment=HOME=/home/admin
ExecStart=/usr/local/bin/setup-terminal-start
Restart=on-failure
RestartSec=3
[Install]
WantedBy=multi-user.target
