# deploy/host/files/etc/systemd/system/project-workbench.service.d/hardening.conf
# GOA security hardening: strip the admin NOPASSWD:ALL sudo grant baked into the
# image (Containerfile) on every (re)start, so agents running as 'admin' cannot
# escalate to root. Appends to the existing ExecStartPost list (does not replace it).
[Service]
ExecStartPost=-/usr/local/sbin/pw-harden-sudoers.sh project-workbench
