# systemd/project-workbench.service.d/auth.conf
# Project Workbench auth enforcement.
#
# Set PW_AUTH_ENFORCE=true to require a valid app session cookie on every
# dashboard request and on every nginx auth_request subrequest. Anonymous
# callers will get a 401 / login redirect.
#
# When false (the safe rollout default), the dashboard treats anonymous
# requests as an implicit admin so existing browser tabs / Claude sessions
# keep working until you have at least one admin user created via
# `sudo pw-user add ... --role admin` and have logged in via /login.
#
# Recommended sequence to enable:
#   1) sudo pw-user add james --role admin --projects '*'
#   2) Open http://<workbench>/login in a browser; verify the cookie works.
#   3) Edit this file: change `false` to `true`.
#   4) sudo systemctl daemon-reload && sudo systemctl restart project-workbench.service
#
# If anything goes wrong, set back to `false` and restart.
[Service]
Environment=PW_AUTH_ENFORCE=false
